# systemd/gloopieguardian.service
[Unit]
Description=GloopieGuardian Django application
After=network.target

[Service]
# Replace 'pi' with the system user that should run the app
User=pi
Group=www-data

# Set to your project's directory
WorkingDirectory=/home/pi/GloopieGuardian

# Update PATH to point to your virtualenv's bin directory
Environment="PATH=/home/pi/GloopieGuardian/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=config.settings"
Environment="PYTHONUNBUFFERED=1"

# Use Django's development server via manage.py runserver.
# Note: `runserver` is for development only — not recommended for production.
# Activate virtualenv before starting the server.
ExecStart=/bin/bash -lc 'source /home/pi/GloopieGuardian/venv/bin/activate && exec python /home/pi/GloopieGuardian/manage.py runserver 0.0.0.0:8000'

Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
